# other/audio-lifesaver.service
[Unit]
Description=Audio lifesaver

[Service]
ExecStart=PATH TO BINARY
Restart=always

[Install]
WantedBy=default.target
